# deploy/ai-agent-learning.service
[Unit]
Description=AgentDesk Multi-Agent Ticket System
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/workspace/ai-agent-learning
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/root/workspace/ai-agent-learning
ExecStart=/root/workspace/ai-agent-learning/venv/bin/python main.py --log-level info
Restart=on-failure
RestartSec=3
StandardOutput=append:/var/log/ai-agent-learning.log
StandardError=append:/var/log/ai-agent-learning.log

[Install]
WantedBy=multi-user.target
